# main.cron
48 2 * * * /home/deaise/voal/project/main.py